# ntfy/rootfs/etc/s6-overlay/s6-rc.d/ntfy/run
#!/command/with-contenv bashio

bashio::log.info "Starting..."

mkdir -p /data/lib/ntfy

export NTFY_AUTH_FILE=/data/lib/ntfy/auth.db
export NTFY_CACHE_FILE=/data/lib/ntfy/cache.db
export NTFY_ATTACHMENT_CACHE_DIR=/data/lib/ntfy/attachments

export NTFY_BASE_URL="$(bashio::config 'base_url')"
export NTFY_AUTH_DEFAULT_ACCESS="$(bashio::config 'default_access')"
export NTFY_ENABLE_LOGIN=true
export NTFY_BEHIND_PROXY="$(bashio::config 'behind_proxy')"
export NTFY_LOG_LEVEL=debug
export NTFY_LISTEN_HTTP=":8080"

for auth in "$(bashio::config 'users')"
do
  name=$(echo "$auth" | jq .name | xargs)
  password=$(echo "$auth" | jq .password | xargs)
  topics_rw=$(echo "$auth" | jq .topics_rw)

  if (exit $(ntfy access "$name"))
  then
    NTFY_PASSWORD="$password" ntfy user change-pass "$name"
  else
    NTFY_PASSWORD="$password" ntfy user add "$name"
  fi

  for topic in "$(echo "$topics_rw" | jq -r '.[]')"
  do
    ntfy access "$name" "$topic" rw
  done
done

bashio::log.info "Finished configuration."

env

ntfy serve --base-url https://home.anpingen.de/