# paperless-ngx-ingress/rootfs/etc/s6-overlay/s6-rc.d/redis/run
#!/command/with-contenv bash

echo "Starting redis in addon ...";

redis-server /etc/redis/redis.conf --daemonize no
