# paperless-ngx-ingress/rootfs/etc/s6-overlay/s6-rc.d/extra/run
#!/command/with-contenv bash

echo "Starting paperless-extra in addon ...";

cd /usr/src/paperless-extra

python manage.py runserver 8001 &

