# paperless-ngx-ingress/rootfs/etc/s6-overlay/s6-rc.d/extra/run
#!/command/with-contenv bash

echo "Starting paperless-extra in addon ...";

cd /usr/src/paperless_extra


pip install -r requirements.txt

python manage.py runserver 0.0.0.0:8001 

