# gitea-act-runner/rootfs/etc/s6-overlay/s6-rc.d/act-runner/run
#!/command/with-contenv bashio

bashio::log.info "Starting..."

export GITEA_RUNNER_NAME=$(bashio::config 'runner_name')
export GITEA_INSTANCE_URL=$(bashio::config 'instance_url')
export GITEA_RUNNER_REGISTRATION_TOKEN=$(bashio::config 'registration_token')
export GITEA_RUNNER_LABELS=$(bashio::config 'labels')

bashio::log.info "Finished configuration."

/usr/bin/gitea_act_runner register \
    --instance "${GITEA_INSTANCE_URL}" \
    --name "${GITEA_RUNNER_NAME}" \
    --token "${GITEA_RUNNER_REGISTRATION_TOKEN}" \
    --labels "${GITEA_RUNNER_LABELS}" 
/usr/bin/gitea_act_runner daemon